# dev-util/mottainai-server/files/mottainai-server.openrc
#!/sbin/openrc-run

### CADDY openRC gist. & Systemd unit used as templative guide..[Unit]Description=Mottainai Agent
### OPENRC on Gentoo AWS is a must , try to convert to a full Sabayon box and offten Kerblewy.. for now can easily deligate my crappy box out. 
## need to add system-d-love systemd:gentoo not convert fully till system is updated.. eslect-init.. etc. equo mask systemd. 
## sysv-init in systemd to help smoth over ... transition , for now need openrc file till systemd is built and all. 
## anycase Mavrick Repo COMUNITY Project I FON'T want to Babysit All DAY. 
## anycase Giving upstream  a toch  of control via the agent wala.. get packages ok... etc. 
### Mavrick  Just a few Crossdev packages or rebuilds that mildly differ from Stock-Sabayon IE QEMU with AARCH64/ARCH64_BE  , and or Experamental Arm64 system packages. ie RPI3 & ROCK64 
## moment though aws is on the el-cheep-o free tier so not much in ram, and 750-1 gig blows chuncks hence why i'd rather agent it and foregt it. 
## and or just push or slave dockers.. /push

Trim my coments out on pull... 

description="Mottainai Agent"
description_reload="Reload configuration"

extra_started_commands="reload"

: ${MOTTIAINAI_SERVER_CONF:=/etc/mottainai/mottainai-server.yaml}
: ${MOTTIAINAI_SERVER_PIDFILE:=/var/run//mottainai-agent}
: ${MOTTIAINAI_SERVER_USER:=mottainai-agent}

depend() {
        need net
        need localmount
        use dns
        after firewall
}

start() {
        ebegin "Smottainai-agent"
        start-stop-daemon --wait 1000 --background --start --exec \
                /usr/bin/mottainai-server\
                --user ${MOTTIAINAI_SERVER_USER} \
                --make-pidfile --pidfile ${MOTTIAINAI_SERVER_PIDFILE} \
                -- -conf "${MOTTIAINAI_SERVER_CONF}" -pidfile "${MOTTIAINAI_SERVER_PIDFILE}" -agree -quiet && \
        chown ${MOTTIAINAI_SERVER_USER}:root ${MOTTIAINAI_SERVER_PIDFILE}
        eend $?
}

stop() {
        ebegin "Stopping mottainai-agent"
        start-stop-daemon --wait 5000 --stop --exec \
                /usr/bin/mottainai-server\
                --user ${MOTTIAINAI_SERVER_USER} \
                --pidfile ${MOTTIAINAI_SERVER_PIDFILE} \
                -s SIGQUIT
        eend $?
}

reload() {
        ebegin "Reloading mottainai-serverconfiguration"
        start-stop-daemon --exec \
                /usr/bin/mottainai-server\
                --user ${MOTTIAINAI_SERVER_USER} \
                --pidfile ${MOTTIAINAI_SERVER_PIDFILE} \
                -s SIGUSR1 && \
        chown ${MOTTIAINAI_SERVER_USER}:root ${MOTTIAINAI_SERVER_PIDFILE}
        eend $?
}
